# conf/goanysync-sync.timer.in
[Unit]
Description=Sync @PACKAGE_NAME@ volumes to disk every hour
Documentation=main:@PACKAGE_NAME@
Wants=local-fs.target

[Timer]
OnCalendar=hourly
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
