# conf/goanysync-sync.service.in
[Unit]
Description=Sync @PACKAGE_NAME@ sync volumes periodically to disk

[Service]
Type=oneshot
ExecStart=@bindir@/@PACKAGE_NAME@ sync
